# deploy/polymarket-gap-audit.service
[Unit]
Description=Polymarket WS Gap Audit
After=network-online.target
Wants=network-online.target
# Emit a one-shot death webhook if the scanner itself crashes — must be
# in [Unit], not [Service].
OnFailure=polymarket-websocket-alert.service

[Service]
Type=oneshot
User=polymarket
Group=polymarket
WorkingDirectory=/opt/polymarket

EnvironmentFile=/opt/polymarket/.env

# Scan heartbeats + reconnect log, emit .gap_manifest.json.
# - ``--alert`` is the primary webhook channel for findings worth paging on.
# - ``--diff-mode`` makes alerts and OnFailure fire only when the audit
#   sees something NEW relative to its previous run.  The first run
#   after deploy alerts once (or stays silent if findings predate the
#   recency cap) and records the current state; subsequent runs only
#   fire on drift — which is what an operator actually wants to see.
# - ``--fail-on-gap-ms 21600000`` (6 hours) and ``--fail-on-bursts``
#   promote NEW findings of exceptional size to a non-zero exit.
# - ``--fail-window-hours 48`` caps the recency of findings worth
#   alerting on so the audit doesn't keep firing on a multi-day-old
#   outage that just got recorded post-restart.
ExecStart=/opt/polymarket/.venv/bin/python /opt/polymarket/scripts/audit_ws_gaps.py \
    --data-dir /opt/polymarket/data \
    --log-dir /var/log/polymarket \
    --alert \
    --diff-mode \
    --fail-on-gap-ms 21600000 \
    --fail-on-bursts \
    --fail-window-hours 48

[Install]
WantedBy=multi-user.target
