# deploy/polymarket-gap-audit.service
[Unit]
Description=Polymarket WS Gap Audit
After=network-online.target
Wants=network-online.target
# Emit a one-shot death webhook if the scanner itself crashes — must be
# in [Unit], not [Service].
OnFailure=polymarket-websocket-alert.service

[Service]
Type=oneshot
User=polymarket
Group=polymarket
WorkingDirectory=/opt/polymarket

EnvironmentFile=/opt/polymarket/.env

# Scan heartbeats + reconnect log, emit .gap_manifest.json.
# - ``--alert`` is the primary webhook on ANY findings (gaps or bursts);
#   that's how operators see sustained degradation.
# - ``--fail-on-gap-ms 21600000`` (6 hours) and ``--fail-on-bursts``
#   promote ONLY exceptional findings to a non-zero exit, so systemd's
#   ``OnFailure=`` fires for major outages and reconnect storms but
#   stays quiet on routine restart-induced gaps (which used to spam
#   ``polymarket-gap-audit.service: Failed`` on every run).
ExecStart=/opt/polymarket/.venv/bin/python /opt/polymarket/scripts/audit_ws_gaps.py \
    --data-dir /opt/polymarket/data \
    --log-dir /var/log/polymarket \
    --alert \
    --fail-on-gap-ms 21600000 \
    --fail-on-bursts

[Install]
WantedBy=multi-user.target
