# deploy/polymarket-historical.service
[Unit]
Description=Polymarket Historical Data Sync + Hugging Face Upload (one-shot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=no
User=polymarket
Group=polymarket
WorkingDirectory=/opt/polymarket

EnvironmentFile=/opt/polymarket/.env

# Incremental scan (checkpoint-based; typically 1-2 pages after the first run)
# then uploads all Parquet data to Hugging Face Hub.
ExecStart=/opt/polymarket/.venv/bin/python -m polymarket_pipeline \
    --historical-only \
    --upload

# Allow unlimited time for the initial massive backfill run
TimeoutStartSec=0
