# deploy/polymarket-restart.service
[Unit]
Description=Restart Polymarket WebSocket service to pick up new markets

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart polymarket-websocket.service
