# deploy/polymarket-restart.service
[Unit]
Description=Restart Polymarket Live service to pick up new markets

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart polymarket-live.service
